# test/script/bundle/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:vllm1]
command=/bin/sh -c "until curl -s http://localhost:8002/health; do echo 'Waiting for vllm2 health...'; sleep 10; done; python3 -m vllm.entrypoints.openai.api_server --model Qwen/Qwen2.5-1.5B --gpu-memory-utilization 0.70 --enable-sleep-mode --max-model-len 1000 --port 8001"
autorestart=false
# Redirect logs to the container's stdout/stderr
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vllm2]
command=python3 -m vllm.entrypoints.openai.api_server --model Qwen/Qwen3-Reranker-0.6B --gpu-memory-utilization 0.2 --max-model-len 1200 --runner pooling --port 8002
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
